# lib/systemd/system/makerbase-timezone-monitor.service
[Unit]
Description=Timezone Monitor Service
After=network.target

[Service]
ExecStart=/root/set-timezone.sh

[Install]
WantedBy=multi-user.target
