# lib/systemd/system/makerbase-beep.service
[Unit]
Description=Touch beep Service
After=KlipperScreen.service

[Service]
ExecStart=/root/beep.sh

[Install]
WantedBy=multi-user.target
